# configs/calsync.service
[Unit]
Description=CalBar Calendar Sync Daemon
Documentation=https://github.com/cpuguy83/calbar

[Service]
Type=simple
ExecStart=%h/.local/bin/calsync
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
